# infrastructure/systemd/arifos-mcp.service
[Unit]
Description=arifOS Constitutional MCP Server
Documentation=https://arifos.arif-fazil.com
After=network.target

[Service]
Type=simple
User=arifos
Group=arifos
WorkingDirectory=/opt/arifos

# Environment
Environment="PATH=/opt/arifos/.venv/bin:/usr/local/bin:/usr/bin"
Environment="ARIFOS_GOVERNANCE_SECRET_FILE=/opt/arifos/secrets/governance.secret"
Environment="ARIFOS_TIMEZONE=Asia/Kuala_Lumpur"
Environment="AAA_MCP_TRANSPORT=http"
Environment="HOST=0.0.0.0"
Environment="PORT=8080"

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/arifos/VAULT999
ReadOnlyPaths=/opt/arifos/secrets

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Main process
ExecStart=/opt/arifos/.venv/bin/python -m arifosmcp.runtime http

# Health check
ExecStartPost=/bin/sh -c 'sleep 5 && curl -f http://localhost:8080/health || exit 1'

[Install]
WantedBy=multi-user.target
